# iwsn-portal-deb/src/deb/control/postinst
#!/bin/sh

# add daemon to runlevel specified in its header
update-rc.d  tr.iwsn-portal defaults || true